# rosbag_recorder/rosbag-recorder.service
[Unit]
Description=ROS2 Bag Recorder
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/
ExecStart=/bin/bash ~/mcs_bridge/rosbag_recorder/record_rosbag.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target